# service/bin/run
#!/bin/bash

config_dir="/service/config"
external_ip_file="${config_dir}/external_ip"

test -d ${config_dir} || mkdir -p ${config_dir}

if [ ! -f "${external_ip_file}" ]; then
  if [ -z "${EXTERNAL_IP}" ]; then
      curl -s http://icanhazip.com > ${external_ip_file}
  else
      echo ${EXTERNAL_IP} > ${external_ip_file}
  fi
fi

redis_conn_string="hostaddr=${REDIS_HOST} port=${REDIS_PORT} password=${REDIS_PASS}"

exec /usr/local/bin/turnserver -n --log-file stdout --mongo-userdb="${COTURN_MONGO_URL}" --redis-statsdb="${redis_conn_string}" --external-ip `cat ${external_ip_file}`
